# deploy/ngrok-tunnel.service
[Unit]
Description=ngrok tunnel for avinash-agent-backend (stable free dev domain)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=avi
ExecStart=/home/avi/.local/bin/ngrok http 8000 --url https://tarantula-prude-cotton.ngrok-free.dev
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
